# auto-kj.service
[Unit]
Description=auto-kj Karaoke Machine
After=network-online.target sound.target graphical.target
Wants=network-online.target

[Service]
Type=simple
User=gabe
Group=gabe
WorkingDirectory=/home/gabe/auto-kj
ExecStart=/home/gabe/auto-kj/.venv/bin/python /home/gabe/auto-kj/auto-kj/main.py
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0
Environment=PULSE_SERVER=unix:/run/user/1000/pulse/native
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus

[Install]
WantedBy=graphical.target
