# auto-kj.service
[Unit]
Description=auto-kj Karaoke Machine
After=sound.target
Wants=sound.target

[Service]
Type=simple
User=gabe
Group=gabe
SupplementaryGroups=audio input render video
TTYPath=/dev/tty7
StandardInput=tty-force
StandardOutput=journal
StandardError=journal
WorkingDirectory=/home/gabe/auto-kj
ExecStart=/home/gabe/auto-kj/.venv/bin/python /home/gabe/auto-kj/auto-kj/main.py
Restart=on-failure
RestartSec=5
LimitRTPRIO=95
LimitMEMLOCK=infinity
Environment=HOME=/home/gabe
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/home/gabe/.env

[Install]
WantedBy=multi-user.target
